# debian/ubuntu-push-client.conf
description "ubuntu push notification client-side daemon"

start on started dbus
stop on stopped dbus

exec /usr/lib/ubuntu-push-client/ubuntu-push-client
