# debian/ubuntu-push-client.conf
description		"Starts the ubuntu push notifications client side daemon"

start on dbus
stop on runlevel [06]

exec /usr/lib/ubuntu-push-client/ubuntu-push-client
